# system/etc/systemd/system/systemd-suspend.service.d/onfailure.conf
# Fire a desktop notification if suspend fails (e.g. an NFS-blocked receivers
# job in D-state prevents the freezer from completing). See
# /usr/local/bin/suspend-failed-notify.
[Unit]
OnFailure=suspend-failed-notify.service
