# systemd/.config/systemd/user/bgovault-sync.timer
[Unit]
Description=Hourly push of bgovault (Obsidian vault)
Documentation=man:systemd.timer(5)

[Timer]
# Push the vault every hour (push-only, safe to run often)
OnBootSec=5min
OnUnitActiveSec=1h

# Anchor to the top of every hour too, so cadence self-heals after sleep/reboot
OnCalendar=*-*-* *:00:00

# Randomize start time by up to 3 minutes to avoid conflicts
RandomizedDelaySec=3min

# Persist timer state across reboots
Persistent=true

[Install]
WantedBy=timers.target
