# systemd/.config/systemd/user/claude-imports.service
[Unit]
Description=Watch Downloads for vault notes and Claude exports

[Service]
ExecStart=%h/notes/bgovault/.scripts/watch-claude-imports.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
