# systemd/.config/systemd/user/password-store-sync.service
[Unit]
Description=Password Store and Dotfiles Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-sync
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
# PrivateTmp=true changes mount namespace in a way that causes OpenSSH 10.2 to reject
# /etc/ssh/ssh_config.d/20-systemd-ssh-proxy.conf (symlink ownership check fails in the
# new namespace). Force git to use only the user SSH config to bypass /etc/ssh entirely.
Environment="GIT_SSH_COMMAND=ssh -F %h/.ssh/config"

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
