# systemd/.config/systemd/user/password-store-sync.timer
[Unit]
Description=Periodic Password Store and Dotfiles Sync
Documentation=man:systemd.timer(5)

[Timer]
# Run sync every 4 hours
OnBootSec=10min
OnUnitActiveSec=4h

# Run on calendar schedule as well (daily at 9am, 1pm, 5pm, 9pm)
OnCalendar=*-*-* 09,13,17,21:00:00

# Randomize start time by up to 5 minutes to avoid conflicts
RandomizedDelaySec=5min

# Persist timer state across reboots
Persistent=true

[Install]
WantedBy=timers.target
