# systemd/.config/systemd/user/password-store-sync.service
[Unit]
Description=Sync dotfiles, claude-private, password-store, and bgovault
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-sync
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Security hardening
# NOTE: PrivateTmp=true is intentionally NOT set. It creates a new user namespace
# (not just a mount namespace), and Linux 6.17 blocks write() on TCP sockets from
# unprivileged user namespaces — which breaks OpenSSH (uses write() not sendto()).
# NoNewPrivileges=true is safe; confirmed not to affect network access.
NoNewPrivileges=true
