# systemd/.config/systemd/user/password-store-sync.timer
[Unit]
Description=Periodic sync of dotfiles, claude-private, password-store, bgovault
Documentation=man:systemd.timer(5)

[Timer]
# Run sync every hour
OnBootSec=5min
OnUnitActiveSec=1h

# Also anchor to top-of-hour wall-clock schedule
OnCalendar=*-*-* *:00:00

# Randomize start time by up to 5 minutes to avoid conflicts
RandomizedDelaySec=5min

# Persist timer state across reboots
Persistent=true

[Install]
WantedBy=timers.target
