# systemd/.config/systemd/user/password-store-sync.service
[Unit]
Description=Password Store and Dotfiles Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-sync
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
