# systemd/.config/systemd/user/spotify-notify.service
[Unit]
Description=Spotify track change notifications
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/spotify-notify
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=default.target
