# system/usr/lib/systemd/system/systemd-suspend.service.d/nvidia-suspend-nofreeze.conf
# nvidia_modeset 595-open has a jump_label bug triggered by freeze_processes().
#
# During suspend, systemd calls freeze_processes() which patches jump labels in
# all loaded kernel modules. nvidia_modeset has an inconsistent jump label at
# nvkms_kthread_q_callback that triggers:
#   kernel BUG at arch/x86/kernel/jump_label.c:73!
# This causes an immediate kernel panic on every suspend attempt.
#
# Setting SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=false bypasses the freeze step,
# avoiding the buggy code path. This is NOT the Xorg/VT-switch workaround — it
# is a workaround for a real nvidia_modeset driver bug with kernel 7.0.0-14.
# Remove this once nvidia-driver is updated and confirmed fixed.
#
# Confirmed crash: /var/crash/202604282111/dmesg.202604282111 (2026-04-28)
#
# Deploy (not stow-managed):
#   sudo cp this-file /usr/lib/systemd/system/systemd-suspend.service.d/
#   sudo systemctl daemon-reload
[Service]
Environment = SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=false
